# crontab
1 * * * * curl -sS https://seahorse-app-mewkl.ondigitalocean.app/websocket/start_connection.json >/proc/1/fd/1 2>/proc/1/fd/2

# An empty line is required at the end of this file for a valid cron file.
